# systemd/llm-router.service
[Unit]
Description=Multi-Backend LLM Router v4.0.0
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/llm-router
Environment="ROUTER_CONFIG=/opt/llm-router/config.json"
Environment="TABBY_TOKENS_PATH=/opt/TabbyAPI/api_tokens.yml"
Environment="TABBY_CONFIG_PATH=/opt/TabbyAPI/config.yml"
Environment="TABBY_MODEL_DIR=/opt/models"
ExecStart=/opt/llm-router/venv/bin/python /opt/llm-router/router.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
