# systemd/tabbyapi.service
[Unit]
Description=TabbyAPI Server - EXL2 Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ivan/TabbyAPI
ExecStart=/home/ivan/TabbyAPI/venv/bin/python main.py
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
